# overlay.d/10-fcos/usr/lib/systemd/system/maintenance-window.service
[Unit]
Description=Upgrade and reboot

[Service]
Type=oneshot
ExecStart=/usr/bin/rpm-ostree upgrade --reboot